# aide/systemd/aide-update.timer.template
[Unit]
Description=AIDE Database Update Timer
Documentation=https://github.com/fidpa/ubuntu-server-security

[Timer]
# Run daily at 04:00
OnCalendar=daily

# Run immediately if missed (e.g., system was off)
Persistent=true

# Optional: Add random delay to prevent load spikes
# RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
